# bsp/arch-dm8168/apps/base/etc/inittab
::sysinit:/etc/rcS
::restart:/sbin/init
ttyO2::askfirst:-/bin/sh
